# cortex/00-System/Services/com.cortex.learning.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cortex.learning.service</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/simonjanke/Projects/cortex/.venv/bin/python3</string>
        <string>/Users/simonjanke/Projects/cortex/00-System/Services/cortex_learner.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/simonjanke/Projects/cortex</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/simonjanke/Projects/cortex/00-System/Services/logs/service_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/simonjanke/Projects/cortex/00-System/Services/logs/service_stderr.log</string>
    
    <key>ThrottleInterval</key>
    <integer>30</integer>
    
    <key>StartInterval</key>
    <integer>1800</integer>
    
    <key>Nice</key>
    <integer>10</integer>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/simonjanke/Projects/cortex/00-System/Algorithms</string>
    </dict>
</dict>
</plist>